# bot/ops/org.wikilean.poll.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>org.wikilean.poll</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/jack/Desktop/LEAN/WikiLean/bot/ops/poll-once.sh</string>
  </array>
  <key>StartInterval</key>
  <integer>1800</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>ProcessType</key>
  <string>Background</string>
  <key>StandardOutPath</key>
  <string>/Users/jack/Library/Logs/WikiLean/poll.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/jack/Library/Logs/WikiLean/poll.err.log</string>
</dict>
</plist>
